# data/me.nalux.Resources.service.in
[Unit]
Description=Resources Daemon

[Service]
Type=dbus
Restart=always
BusName=me.nalux.Resources
ExecStart=@bindir@/resources-daemon

[Install]
WantedBy=multi-user.target